# config/openclaw-gateway.service.template
[Unit]
Description=OpenClaw Gateway
After=network.target

[Service]
Type=simple
ExecStart=__NODE_BIN__ __OPENCLAW_BIN__
Restart=on-failure
RestartSec=5
Environment=HOME=__HOME__
Environment=PATH=__NODE_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=default.target
